# .config/systemd/user/move-notification.timer
[Unit]
Description=Send notification to move every 20 min

[Timer]
OnActiveSec=20min
OnUnitInactiveSec=20min

[Install]
WantedBy=timers.target
