# .config/systemd/user/move-notification.service
[Unit]
Description=Send notification to move

[Service]
ExecStart=/usr/bin/notify-send -u critical 'Beweeg!' 'Tijd om je rug wat te bewegen!'
